# deepcool-ak-series-digital.service
[Unit]
Description=Displaying the current temperature and utilization on the Deepcool AK Series CPU cooler

[Service]
Type=simple
WorkingDir=~
ExecStart= /usr/bin/python .local/bin/deepcool-ak-series-digital.py
Restart=always

[Install]
WantedBy=default.target
